# payload/Library/LaunchAgents/com.anothertoolappleshouldhaveprovided.nudge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.anothertoolappleshouldhaveprovided.nudge</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Application Support/nudge/Resources/nudge</string>
		<!-- Cutoff Date -->
		<!--<string>--cutoffdate</string> -->
		<!--<string>2018-12-31-00:00</string> -->
		<!-- -->
		<!-- Logos -->
		<!--<string>--logopath</string> -->
		<!--<string>/Some/Custom/Path/company_logo.png</string> -->
		<!--<string>--updatesspath</string> -->
		<!--<string>/Some/Custom/Path/update_ss.png</string> -->
		<!-- -->
		<!-- OS -->
		<!--<string>--minimumosversion</string> -->
		<!--<string>10.13</string> -->
		<!--<string>--pathtoapp</string> -->
		<!--<string>/Applications/Install macOS High Sierra.app</string> -->
		<!-- -->
		<!-- Text-->
		<!--<string>--updatetext</string> -->
		<!--<string>A security update is required on your machine.</string> -->
		<!--<string>--h1text</string> -->
		<!--<string>Ready to start the update?</string> -->
		<!--<string>--h2text</string> -->
		<!--<string>Click on the button below.</string> -->
		<!--<string>--paragraph1</string> -->
		<!--<string>A fully up-to-date device is required to ensure that IT can your accurately protect your computer.</string> -->
		<!--<string>--paragraph2</string> -->
		<!--<string>If you do not update your computer, you may lose access to some items necessary for your day-to-day tasks.</string> -->
		<!--<string>--paragraph3</string> -->
		<!--<string>To begin the update, simply click on the button below and follow the provided steps.</string> -->
		<!--<string>--subtitletext</string> -->
		<!--<string>A friendly reminder from your local IT team</string> -->
		<!--<string>--titletext</string> -->
		<!--<string>macOS Upgrade</string> -->
		<!-- -->
		<!-- Timers -->
		<!--<string>--cutoffdatewarning</string> -->
		<!--<string>3</string> -->
		<!--<string>--timerday1</string> -->
		<!--<string>600</string> -->
		<!--<string>--timerday3</string> -->
		<!--<string>7200</string> -->
		<!--<string>--timerelapsed</string> -->
		<!--<string>10</string> -->
		<!--<string>--timerfinal</string> -->
		<!--<string>60</string> -->
		<!--<string>--timerinital</string> -->
		<!--<string>14400</string> -->
		<!-- Timers Boolean-->
		<!--<string>--notimer</string> -->
		<!-- URLs -->
		<!--<string>--moreinfourl</string> -->
		<!--<string>https://google.com</string> -->
		<!-- -->
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Library/Application Support/nudge/Logs/nudge.log</string>
	<key>StandardOutPath</key>
	<string>/Library/Application Support/nudge/Logs/nudge.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
	</array>
</dict>
</plist>
