# payload/Library/LaunchAgents/com.erikng.nudge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.erikng.nudge</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/ManagedFrameworks/Python/Python3.framework/Versions/Current/bin/python3</string>
		<string>/Library/nudge/Resources/nudge</string>
		<string>--jsonurl=https://fake.domain.com/path/to/config.json</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Library/nudge/Logs/nudge.log</string>
	<key>StandardErrorPath</key>
	<string>/Library/nudge/Logs/nudge.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
	</array>
</dict>
</plist>
